# V1/inittab
::sysinit:/etc/init.d/rcS

::askfirst:-/bin/ash
::sysinit:/bin/mount -t proc proc /proc

::restart:/sbin/init

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a